# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:admin]
command=python run.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

# If your website server is started automatically by run.py when the plugin is enabled,
# you do NOT need a separate [program:website] entry.
# If you ever split the website server into a separate script, add:
; [program:website]
; command=python website_server.py
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stderr_logfile=/dev/stderr
